# docker/efergy-poller/crontab
* * * * * root /usr/src/app/project_env.sh; /usr/local/bin/python3 /usr/src/app/efergy-poller.py >> /var/log/cron.log 2>&1
